# docker/etc/services.d/cockpit/run
#!/usr/bin/env bash

set -eu
set -o pipefail


export PATH="$PATH:/usr/libexec"


exec cockpit-ws --no-tls


