# root/etc/services.d/novnc/run
#!/usr/bin/env bash

set -eu
set -o pipefail


exec /novnc/utils/launch.sh --listen 66080 --vnc localhost:65900

