# root/etc/services.d/domain-start/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail


s6-svwait -U /var/run/s6/services/nat-init


VMDK="/config"
DOMAIN_XML="$VMDK/$VM_NAME.xml"


check_running() {
  local running="$(virsh list --name --state-running | tr -d '\n' | wc -l)"
  if [[ "$running" -eq 0 ]]
  then
    exit
  fi
}


if [[ -f "$DOMAIN_XML" ]]
then
  printf '%s\n' '--------------------------------------------------------'
  cat "$DOMAIN_XML"
  printf '%s\n' '--------------------------------------------------------'

  virsh define "$DOMAIN_XML"
  virsh start "$VM_NAME"

  while true
  do
    sleep 3
    check_running
  done
else
  >&2 printf '%s\n' '--------------------------------------------------------'
  >&2 printf '%s\n' 'WARNING -- VM Not Found'
  >&2 printf '%s\n' 'Build VM using -- CMD new'
  >&2 printf '%s\n' 'See https://github.com/ms-jpq/windows-in-docker'
  >&2 printf '%s\n' '--------------------------------------------------------'
  sleep infinity
fi

