# base/_root/etc/services.d/domain-start/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

s6-svwait -U /var/run/s6/services/net-define


sleep infinity
