# base/_root/etc/services.d/domain-start/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

s6-svwait -U /var/run/s6/services/net-define

VMDK="/config"
DOMAIN_XML="$VMDK/$VM_NAME.xml"

if [[ -f "$DOMAIN_XML" ]]
then
  virsh define "$DOMAIN_XML"
  while true
  do
    virsh start "$VM_NAME" || true
    sleep 3
  done
else
  hr
  >&2 echo WARNING -- VM Not Found
  >&2 echo Build VM using 'CMD :: new'
  hr
  sleep infinity
fi

