# base/_root/etc/services.d/kvm/run
#!/usr/bin/env bash

set -eu
set -o pipefail

cd /vmrc || exit 1

s6-svwait /var/run/s6/services/libvirtd

while true
do
  for rc in **
  do
    ."/$rc"
  done
  sleep 3
done
